# systemd/system/manx.service
[Unit]
Description="DNS blocklist merger, parser, inspector"

[Service]
ExecStart=/usr/bin/manx-daemon