# systemd/system/manx.timer
[Unit]
Description="Run manx every three days"

[Timer]
OnActiveSec=5min
OnUnitActiveSec=72h
Unit=manx.Service

[Install]
WantedBy=multi-user.target